# backup.service
[Unit]
Description=Backup diário do banco Estoque

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/estoque
ExecStart=/opt/estoque/backup_prod.sh
