# backup.timer
[Unit]
Description=Timer para backup diário do banco Estoque

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
